# alpine-minirootfs/etc/profile.d/motd.sh
echo "Weclome to One File Linux!
To report a bugs: https://github.com/zhovner/OneFileLinux

To change console font size: fontsize [1-9]
To connect WiFi read: /usr/share/doc/wifi-howto.txt

"
