# alpine-minirootfs/etc/profile.d/motd.sh
echo "Weclome to One File Linux!

Feel free to report a bugs: https://github.com/zhovner/OneFileLinux

You can change console font size: fontsize [1-9]

"
