# alpine-minirootfs/etc/init.d/acpid
#!/sbin/openrc-run

name="busybox acpid"
command="/sbin/acpid"
command_args="$ACPID_OPTS"
pidfile="/var/run/acpid.pid"
start_stop_daemon_args="-g wheel -k 027"

depend() {
	need dev localmount
	after hwdrivers modules
	keyword -vserver -lxc
}

