# alpine-minirootfs/etc/init.d/crond
#!/sbin/openrc-run

name="busybox $SVCNAME"
command="/usr/sbin/$SVCNAME"
pidfile="/var/run/$SVCNAME.pid"
command_args="$CRON_OPTS"

depend() {
	need localmount
	need logger
}
