# alpine-minirootfs/etc/init.d/httpd
#!/sbin/openrc-run

name="busybox $SVCNAME"
command="/usr/sbin/$SVCNAME"
command_args="$HTTPD_OPTS"

depend() {
	need net localmount
	after firewall
}

