# alpine-minirootfs/etc/init.d/klogd
#!/sbin/openrc-run

name="busybox klogd"
command="/sbin/klogd"
command_args="${KLOGD_OPTS}"
pidfile="/var/run/klogd.pid"

depend() {
	need clock hostname localmount
	before net
	keyword -vserver -lxc
}
