# alpine-minirootfs/etc/init.d/rdate
#!/sbin/openrc-run

depend() {
	need net
	after firewall
	keyword -vserver -lxc
}

start() {
	if [ -z "$NTP_HOST" ] ; then
		eerror "Please set NTP_HOST in /etc/conf.d/rdate"
		return 1
	fi
	ebegin "Running busybox rdate"
	rdate "$NTP_HOST"
	eend $?
}

